# systemd/mcga.service
[Unit]
Description=update mcga

[Service]
Type=oneshot
ExecStart=/usr/bin/mcga.sh

[Install]
WantedBy=multi-user.target
